# cron.d/restic-kopano
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

0 */6 * * * root /root/kopano-docker-restic-backup/backup-job
